# fastnet2ip.service
# sudo cp fastnet2ip.service /etc/systemd/system/
# sudo systemctl daemon-reload
# sudo systemctl enable fastnet2ip
# sudo systemctl start fastnet2ip
#
# sudo systemctl status fastnet2ip
# sudo systemctl stop fastnet2ip
# sudo systemctl restart fastnet2ip
# journalctl -u fastnet2ip -f

[Unit]
Description=fastnet2ip Service
####################################### CHANGE THIS IF PORT CHANGES
After=dev-ttySTM3.device
BindsTo=dev-ttySTM3.device

[Service]
Type=simple
User=root
####################################### CHANGE THIS PATH to your install
# Recommended: a dedicated venv, e.g.  python3 -m venv /opt/fastnet2ip && /opt/fastnet2ip/bin/pip install fastnet2ip
WorkingDirectory=/opt/fastnet2ip
# NMEA 2000 output (uncomment one ExecStart only):
ExecStart=/opt/fastnet2ip/bin/fastnet2ip --output nmea2000 --serial /dev/ttySTM3 --udp-port 2000 --n2k-format ydwg --n2k-src 201 --n2k-pri 4 --ignore-gps --log-level INFO
# NMEA 0183 output:
#ExecStart=/opt/fastnet2ip/bin/fastnet2ip --output nmea0183 --serial /dev/ttySTM3 --udp-port 2002 --log-level INFO
Restart=always
RestartSec=10

# === RESOURCE LIMITS ===
OOMScoreAdjust=-700
OOMPolicy=continue
MemoryMax=128M
MemoryHigh=96M
TimeoutStopSec=30


# === LOGGING ===
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fastnet2ip
Environment=PYTHONUNBUFFERED=1


[Install]
WantedBy=multi-user.target
